# dashboard/bin/calendar.conf
# calendar

description "calendar screen-writing daemon"
author "Seb Bacon seb.bacon@gmail.com>"

start on started lab126_gui
stop on stopping lab126

respawn
# Maximum of ten times, with gap of 60 seconds
respawn limit 10 60

exec start-stop-daemon -S -b -m -p /tmp/calendar.pid -x "/bin/sh"  -- /mnt/base-us/extensions/dashboard/bin/start.sh